# crontab
# for Notion Scripts
# execute once.sh everyday on 5:00 am
0 5 * * * sh /root/Scripts/NotionScripts/once.sh
